# service/haproxy/run
#!/usr/bin/env bash

echo "=> Starting HAProxy ..."

exec 2>&1
exec /usr/sbin/haproxy -f /etc/haproxy/haproxy.cfg
